# services/launchd/com.andersbakken.rtags.agent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.andersbakken.rtags.agent</string>
    <key>ProgramArguments</key>
    <array>
      <string>/bin/sh</string>
      <string>-c</string>
      <string>/Users/tsc/opt/local/bin/rdm --launchd --inactivity-timeout 300 --log-file /Users/tsc/Library/Logs/rtags.launchd.log</string>
    </array>
    <key>Sockets</key>
    <dict>
      <key>Listener</key>
      <dict>
        <key>SockPathName</key>
        <string>/Users/tsc/.rdm</string>
      </dict>
    </dict>
  </dict>
</plist>
